# supervisord.conf
[supervisord]
nodaemon=true
; Pipe supervisord's own log to docker stdout so everything appears in `docker logs`
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
command=node /app/backend/server.js
directory=/app/backend
autostart=true
autorestart=true
; Force PORT=5000 — prevents host/platform PORT injection from overriding this
environment=PORT="5000",BACKEND_PORT="5000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5
startretries=3

[program:frontend]
command=node /app/frontend/server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOSTNAME="0.0.0.0",PORT="3000"
startsecs=5
startretries=3
