# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=arwapark

[program:backend]
command=node /app/backend/server.js
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.err
environment=NODE_ENV="%(ENV_NODE_ENV)s",PORT="5000"
startsecs=5
startretries=3

[program:frontend]
command=node /app/frontend/server.js
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend.err
environment=NODE_ENV="%(ENV_NODE_ENV)s",PORT="3000",HOSTNAME="0.0.0.0"
startsecs=5
startretries=3
